# script/upstart/redis.conf
start on (net-device-up and local-filesystems)
stop on runlevel [016]

script
  exec redis-server
end script

