# script/upstart/resque.conf
start on started redis
stop on stopping redis

script
  exec sudo -u pophealth -i 'rvm use 1.9.2; cd /home/pophealth/popHealth ; ./start_worker.sh' >> /tmp/upstart.log 2>&1
end script
