# script/upstart/resque.conf
start on started redis
stop on stopping redis

script
  exec sudo -u pophealth /home/pophealth/popHealth/script/start_resque.sh >> /tmp/resque.log 2>&1
end script
