# script/sharding/templates/upstart.conf.erb
# Ubuntu upstart file at /etc/init/<%= name %>.conf

limit nofile 20000 20000

kill timeout 300 # wait 300s between SIGTERM and SIGKILL.

pre-start script
    mkdir -p /var/lib/<%= name %>/
    mkdir -p /var/log/<%= name %>/
    touch /var/log/<%= name %>/mongodb.log

    chown -R <%= user || "mogodb" %> /var/lib/<%= name %>/
    chown -R mongodb /var/log/<%= name %>/
end script

start on runlevel [2345]
stop on runlevel [06]

script

  exec sudo -u mongodb /usr/bin/mongod <%= options %> /etc/<%= name %>.conf

end script